# ansible/roles/icecube/templates/systemd/user/randomize_pvpn.timer.j2
[Unit]
Description=Schedule random VPN exit node choice

[Timer]
OnUnitActiveSec={{ protonvpn.peer_rotation }}
Persistent=true

[Install]
WantedBy=default.target
